# supervisord.conf
# supervisor config file, for mac

[unix_http_server]
file=/tmp/supervisor.sock   ; (the path to the socket file)

[supervisord]
logfile=/tmp/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=/tmp/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket

[program:rabbitmq]
command=/usr/local/sbin/rabbitmq-server
numprocs=1
autostart=false
autorestart=true
stderr_logfile=NONE

[program:redis-server]
command=redis-server
numprocs=1                    ; number of processes copies to start (def 1)
autostart=false                ; start at supervisord start (default: true)
autorestart=true        ; whether/when to restart (default: unexpected)
stderr_logfile=NONE

[program:mongod]
command=/usr/local/bin/mongod --config /usr/local/etc/mongod.conf
autostart=false                ; start at supervisord start (default: true)
